# config/crontab
#! /bin/bash
# crontab -e as ubuntu

# hide the cursor if it hasn't moved
@reboot /usr/bin/unclutter -idle 0.01 -root &

# Turn on the display at 8 AM
0 8 * * * /home/ubuntu/src/power/display_on

# fetch the morning paper
45 8 * * * /home/ubuntu/src/nyt/fetch
# enough paper
0 10 * * * pkill evince

# display art
0 10 * * * /usr/bin/python3 /home/ubuntu/src/art/main.py

# Turn off the display at 10:30 PM
30 20 * * * /home/ubuntu/src/power/display_off
