# config/crontab
# crontab -e as ubuntu

@reboot /home/ubuntu/src/eo1/wakeup

# good morning
0 8 * * * /home/ubuntu/src/eo1/wakeup

# display art
0 10 * * * /home/ubuntu/src/eo1/cli art

# Dim a bit
0 20 * * * /home/ubuntu/src/eo1/cli brightness 5

# Dim more
0 22 * * * /home/ubuntu/src/eo1/cli brightness 3

# Turn off the display at 11 PM
0 23 * * * /home/ubuntu/src/eo1/cli off
